# files/vpn.service
[Unit]
Description="A script to backup all home dirs to usb drive"

[Service]
ExecStart=/usr/bin/curl -4s https://ifconfig.me
